# systemd/make_repo_torrents.timer
[Unit]
Description=Periodic creation of pacman2pacman torrent files

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
