# systemd/update-lastsync.service
[Unit]
Description=Update ${FTP_BASE}/lastsync

[Service]
Type=oneshot
User=repo
ExecStart=/opt/dbscripts/cron-jobs/update-lastsync
